# salt/observer/config/etc-init-metabase.conf
start on (local-filesystems and net-device-up IFACE!=lo)
env MB_JETTY_PORT=8000
setuid elife
chdir /srv/metabase
exec java -jar metabase.jar
